# cli/src/deb/control/postinst
#!/bin/sh

set -e

if [ -x "/etc/init.d/jmx-zabbix" ]; then
	invoke-rc.d jmx-zabbix start || exit $?
fi

