# cli/src/deb/control/prerm
#!/bin/sh
set -e

if [ -x "/etc/init.d/jmx-zabbix" ]; then
	invoke-rc.d jmx-zabbix stop || exit $?
fi
